# scripts/rc.local
#!/bin/sh
echo -e 'power on\ndffault-agent\ndiscoverable on\ntrust 00:17: \t \nquit' | bluetoothctl
